# scripts/pve2-cpu-monitor.service
[Unit]
Description=PVE2 CPU Usage Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/pve2-cpu-monitor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
